# supervisord.conf
[supervisord]


[supervisorctl]


[program:horizon]
command=php artisan horizon
autostart=true

[program:echo-server]
command=laravel-echo-server start
autostart=true

